# source/master/automator/etc/init/xe-automator.conf
# xe automator

description     "XenServer Guest Configuration Automator"

start on stopped rc RUNLEVEL=[2345]

script
        /usr/sbin/xe-set-network
        /usr/sbin/xe-set-hostname
        /usr/sbin/generate-sshd-keys
end script
